# postgresql/postgres.cron
0 0 * * * postgres /usr/bin/vacuumdb -z -U updateservce -p 5439 updateservice >& /dev/null

