# postgresql/postgres.cron
0 0 * * * postgres /usr/bin/vacuumdb -z -U updateservce updateservice >& /dev/null

