# postgresql/postgres.cron.in
0 0 * * * postgres @BINDIR@/vacuumdb -z -U updateservce updateservice >& /dev/null

