# postgresql/postgres.cron.in
0 0 * * * postgres @BINDIR@/vacuumdb -az -U rus >& /dev/null

