# distro/postgresql/postgres.cron
0 0 * * * postgres /usr/bin/vacuumdb -z -q -U updateservice -p 5439 updateservice

