# distro/post-install
#!/bin/bash
chkconfig nginx on
chkconfig gunicorn on
chkconfig pgbouncer on
chkconfig postgresql-updateservice on

chkconfig rmake off
exit 0
